# deploy/conduit-dsgndude.service
[Unit]
Description=Conduit Node (DsgnDude)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/conduit
ExecStart=/root/conduit/target/release/conduit-setup \
  --storage-dir /var/lib/conduit-node \
  --port 9741 \
  --rpc-host 137.184.18.6 --rpc-port 38332 \
  --rpc-user lightning --rpc-password lightning \
  --http-port 3006 \
  --public-ip 147.182.138.182 \
  --registry-url http://157.230.238.79:3003 \
  --dashboard /root/conduit-ui/dashboard.html \
  --alias DsgnDude \
  serve
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=conduit-dsgndude

[Install]
WantedBy=multi-user.target
